# stealth-overlay.service
[Unit]
Description=Stealth Overlay Service
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/stealth-overlay %h/.config/stealth-overlay/overlay.yml
Restart=always
RestartSec=5
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
Environment=LD_PRELOAD=%h/.local/lib/libstealth_hook.so

# Security settings - relaxed for LD_PRELOAD compatibility
PrivateTmp=true

[Install]
WantedBy=default.target
